# crontabrc/crontabrc/silent.cron
0 23,0-10 * * * /Users/huangyingw/loadrc/crontabrc/silent.sh
