# crontabrc/.imacpro.cron
*/10 * * * * /Users/huangyingw/loadrc/crontabrc/set_brightness.sh
