# crontabrc/.imacpro.cron
* * * * * /Users/huangyingw/loadrc/crontabrc/set_brightness.sh
*/10 * * * * /Users/huangyingw/loadrc/bashrc/ports_forward.sh > /Users/huangyingw/loadrc/ports_forward.log 2>&1
