# crontabrc/.imacpro.cron
* * * * * /Users/huangyingw/loadrc/crontabrc/set_brightness.sh
