# crontabrc/.imacpro.cron
0 8 * * * /Users/huangyingw/loadrc/crontabrc/setup_cron.sh
