# crontabrc/.imacpro.cron
* * * * * /Users/huangyingw/loadrc/crontabrc/setup_cron.sh
