# crontabrc/.imacpro.cron
0 8 * * * /root/loadrc/crontabrc/mldonkey_start.sh
