# crontabrc/.mbp.cron
*/10 * * * * /root/loadrc/crontabrc/mldonkey_start.sh
