# crontabrc/.mbp.cron
*/10 * * * * /root/loadrc/dropboxrc/dropbox_start.sh
