# crontabrc/.mini.cron
30 23 * * * /root/loadrc/crontabrc/mldonkey_start.sh
0 7-23 * * * /root/loadrc/crontabrc/mldonkey_stop.sh
