# crontabrc/.mini.cron
30 * * * * /root/loadrc/crontabrc/mldonkey_start.sh
0 * * * * /usr/bin/expressvpn connect smart
0 22-23,0-6 * * * /root/myproject/git/interview/skygragon/leetcode-cli/watch.sh
