# crontabrc/.mini.cron
0 0 * * * root /root/loadrc/crontabrc/mldonkey_start.sh
0 7-23 * * * root /root/loadrc/crontabrc/mldonkey_stop.sh
