# crontabrc/.mldonkey.cron
0 * * * * /usr/bin/expressvpn connect smart
*/1 * * * * /root/loadrc/crontabrc/mldonkey_start.sh
