# crontabrc/.server.cron
* 23,0-7 * * * /root/loadrc/bashrc/sleep.sh
