# crontabrc/.server.cron
* * * * * /root/loadrc/bashrc/watch_jdupes.sh >> /root/loadrc/watch_jdupes.log 2>&1
