# crontabrc/silent.cron
0 22-23,0-8 * * * /Users/huangyingw/loadrc/crontabrc/silent.sh
0 22-23,0-8 * * * /Users/yhuang/loadrc/crontabrc/silent.sh
