# crontabrc/ssh.cron
* * * * * /Users/huangyingw/loadrc/macosrc/restart_ssh.sh
* * * * * /Users/ying.huang/loadrc/macosrc/restart_ssh.sh
