# crontabrc/ssh.cron
*/5 * * * * /Users/huangyingw/loadrc/macosrc/restart_ssh.sh
*/5 * * * * /Users/ying.huang/loadrc/macosrc/restart_ssh.sh
