# crontabrc/ssh.cron
*/10 * * * * /Users/huangyingw/loadrc/macosrc/restart_ssh.sh
*/10 * * * * /Users/ying.huang/loadrc/macosrc/restart_ssh.sh
