# crontabrc/ssh.cron
*/3 * * * * /Users/huangyingw/loadrc/macosrc/restart_ssh.sh
*/3 * * * * /Users/ying.huang/loadrc/macosrc/restart_ssh.sh
