# crontabrc/wifi.cron
*/10 * * * * /root/loadrc/wifirc/enable.sh
