# crontabrc/zerotier.cron
0,10,20,30,40,50 * * * * /root/loadrc/zerotierrc/zerotier_restart.sh
0,10,20,30,40,50 * * * * /Users/huangyingw/loadrc/zerotierrc/zerotier_restart.sh
