# crontabrc/zerotier.cron
*/10 * * * * /root/loadrc/zerotierrc/zerotier_restart.sh
*/10 * * * * /Users/huangyingw/loadrc/zerotierrc/zerotier_restart.sh
