# crontabrc/zerotier.cron
0 * * * * /root/loadrc/zerotierrc/zerotier_restart.sh
