# crontabrc/zerotier.cron
*/10 * * * * /root/loadrc/zerotierrc/restart.sh
*/10 * * * * /Users/huangyingw/loadrc/zerotierrc/restart.sh
