# crontabrc/zerotier.cron
0 * * * * /root/loadrc/zerotierrc/zerotier_restart.sh
0 * * * * /Users/huangyingw/loadrc/zerotierrc/zerotier_restart.sh
